# data/hpm-platform.service
[Unit]
Description=OMEN Command Center — Platform Info Service
After=dbus.service systemd-logind.service
Requires=dbus.service systemd-logind.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /usr/libexec/hp-manager/services/platform_service.py
Restart=on-failure
RestartSec=3

# Security hardening — namespace-creating directives (PrivateTmp, ProtectHome,
# PrivateNetwork) disabled to fix status=226/NAMESPACE on Ubuntu 24.04+/Zorin OS 18.
# Non-namespace directives are kept for defence-in-depth with proper hardware capabilities.
ProtectHome=false
PrivateTmp=false
PrivateNetwork=false
NoNewPrivileges=true
LockPersonality=true
RestrictSUIDSGID=true
ProtectSystem=strict
ReadWritePaths=/sys /proc /etc/udev /etc/hp-manager
ProtectKernelModules=true
ProtectKernelTunables=false
ProtectControlGroups=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
