# data/hpm-platform.service
[Unit]
Description=OMEN Command Center — Platform Info Service
After=dbus.service systemd-logind.service
Requires=dbus.service systemd-logind.service

[Service]
Type=dbus
BusName=com.yyl.hpmanager.platform
ExecStart=/usr/bin/python3 /usr/libexec/hp-manager/services/platform_service.py
Restart=on-failure
RestartSec=3

# Security hardening — namespace-creating directives (PrivateTmp, ProtectHome,
# PrivateNetwork) disabled to fix status=226/NAMESPACE on Ubuntu 24.04+/Zorin OS 18.
# Non-namespace directives are kept for defence-in-depth with proper hardware capabilities.
ProtectHome=false
PrivateTmp=false
PrivateNetwork=false
NoNewPrivileges=true
LockPersonality=true
RestrictSUIDSGID=true
ProtectSystem=strict
ReadWritePaths=/sys /proc /etc/udev /etc/hp-manager /tmp
ProtectKernelModules=true
ProtectKernelTunables=false
ProtectControlGroups=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_NET_ADMIN

# D-Bus activated on demand — no [Install]/WantedBy so it does NOT start at
# boot.  It is launched when the GUI first talks to it and exits itself after
# an idle period (see HPM_IDLE_TIMEOUT in dbus_helpers.py).
